# roles/laptophw/templates/powertop.service.j2
[Unit]
Description=Powertop tunings

[Service]
ExecStart={{ laptophw.powertop.bin }} --auto-tune
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
